# recipes-core/initscripts/initscripts/banner.sh
#!/bin/sh
### BEGIN INIT INFO
# Provides:         banner
# Required-Start:   $all
# Required-Stop:
# Default-Start:    99
# Default-Stop:     99
### END INIT INFO

#BLANK, start of cluster is done in startCluster.sh
